# config/systemd/user/auto_share.service
[Unit]
Description=Check the network mounts

[Service]
Type=simple
ExecStart=/home/charles/scripts/auto_share
